# supervisord.conf
[supervisord]
nodaemon=true

[program:flask]
command=flask run --host=0.0.0.0
directory=/app
environment=FLASK_APP="app.py",FLASK_ENV="production",CONFIG_PATH="/config"
autostart=true
autorestart=true
stderr_logfile=/var/log/flask.err.log
stdout_logfile=/var/log/flask.out.log

[program:cron]
command=cron -f
autostart=true
autorestart=true
stderr_logfile=/var/log/cron.err.log
stdout_logfile=/var/log/cron.out.log
